# etc/cron.hourly/afad-autoupdate
#!/bin/bash
yes | dpkg --configure
apt update && DEBIAN_FRONTEND=noninteractive apt full-upgrade -y